# deploy/d-brain-bot.service
[Unit]
Description=d-brain Telegram Bot
After=network.target

[Service]
Type=simple
User=khabarovru
WorkingDirectory=/home/khabarovru/.openclaw/workspace/skills/agent-second-brain
ExecStart=/home/khabarovru/.local/bin/uv run python -m d_brain
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/khabarovru/.openclaw/workspace/skills/agent-second-brain/.env

[Install]
WantedBy=multi-user.target
