# deploy/d-brain-weekly.service
[Unit]
Description=d-brain Weekly Digest
After=network.target

[Service]
Type=oneshot
User=khabarovru
WorkingDirectory=/home/khabarovru/.openclaw/workspace/skills/agent-second-brain
ExecStart=/home/khabarovru/.local/bin/uv run python scripts/weekly.py
EnvironmentFile=/home/khabarovru/.openclaw/workspace/skills/agent-second-brain/.env
StandardOutput=journal
StandardError=journal
